# sdi/launchd/net.hylograph.sdi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.hylograph.sdi</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/afc/work/afc-work/agent-teams/sdi/launchd/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/afc/work/afc-work/agent-teams/sdi</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/afc/Library/Logs/sdi/router.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/afc/Library/Logs/sdi/router.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/afc</string>
        <key>SDI_HOST</key>
        <string>mbp</string>
        <key>MARGINALIA_URL</key>
        <string>http://andrews-mac-mini:3100</string>
        <key>SDI_FETCH_TIMEOUT_MS</key>
        <string>10000</string>
        <!-- DeepStar-managed music-rig ports (#191). Stop-gap until a
             managed_by column on project_servers makes this declarative. -->
        <key>SDI_EXCLUDE_PORTS</key>
        <string>3012,3060,3061,20808,57120</string>
    </dict>
</dict>
</plist>
