# project-tracker/tools/launchd/net.hylograph.marginalia.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.hylograph.marginalia.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/afc/work/afc-work/agent-teams/project-tracker/tools/launchd/start-api.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/afc/work/afc-work/agent-teams/project-tracker</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/afc/Library/Logs/marginalia/api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/afc/Library/Logs/marginalia/api.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/afc</string>
    </dict>
</dict>
</plist>
